# config/systemd/autovoice.service
[Unit]
Description=AutoVoice Voice Conversion Service
Requires=docker.service
After=docker.service nvidia-persistenced.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/thordrive/autovoice
ExecStartPre=/usr/bin/docker compose -f %h/thordrive/autovoice/docker-compose.yaml pull --quiet backend frontend
ExecStart=/usr/bin/docker compose -f %h/thordrive/autovoice/docker-compose.yaml up -d backend frontend
ExecStop=/usr/bin/docker compose -f %h/thordrive/autovoice/docker-compose.yaml down
ExecReload=/usr/bin/docker compose -f %h/thordrive/autovoice/docker-compose.yaml restart backend frontend
Restart=on-failure
RestartSec=10s
User=kp
Group=kp

# Environment
Environment="COMPOSE_PROJECT_NAME=autovoice"

# Security
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
